# scripts/deployment/systemd/kafka-pipeline-xact-upload@.service
[Unit]
Description=Kafka Pipeline XACT Upload Worker (Instance %i)
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=kafka-pipeline
Group=kafka-pipeline
WorkingDirectory=/opt/kafka-pipeline/src

# Instance-specific worker ID
Environment="WORKER_ID=xact-upload-%i"

# Load environment from shared config
EnvironmentFile=/opt/kafka-pipeline/config/.env

# Python executable and worker command
ExecStart=/opt/kafka-pipeline/venv/bin/python -m kafka_pipeline \
    --worker xact-upload \
    --metrics-port 801%i \
    --log-level INFO

# Restart policy
Restart=always
RestartSec=10

# Resource limits
MemoryMax=2G
CPUQuota=150%

# Security hardening
NoNewPrivileges=true
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=kafka-pipeline-xact-upload-%i

[Install]
WantedBy=multi-user.target
