# scripts/deployment/systemd/kafka-pipeline-claimx-delta-writer.service
[Unit]
Description=Kafka Pipeline ClaimX Delta Writer
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=kafka-pipeline
Group=kafka-pipeline
WorkingDirectory=/opt/kafka-pipeline/src

# Worker ID
Environment="WORKER_ID=claimx-delta-writer"

# Load environment from shared config
EnvironmentFile=/opt/kafka-pipeline/config/.env

# Python executable and worker command
ExecStart=/opt/kafka-pipeline/venv/bin/python -m kafka_pipeline \
    --worker claimx-delta-writer \
    --metrics-port 8202 \
    --log-level INFO

# Restart policy
Restart=always
RestartSec=10

# Resource limits
MemoryMax=1G
CPUQuota=100%

# Security hardening
NoNewPrivileges=true
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=kafka-pipeline-claimx-delta-writer

[Install]
WantedBy=multi-user.target
